# etc/service/upsource/run
#!/bin/sh

# Redirect stderr so everything ends up in the log file
exec 2>&1
exec $UPSOURCE_HOME_DIR/Upsource/bin/upsource.sh run 2>&1
